# native/jenkins/upstart/centos6-c++-build-slave.conf
description "Centos 6 C++ build slave container"
author "Sam Meder"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a 8f3d556861f7
